# com.edith.claude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.edith.claude</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/randywilson/Desktop/edith-v3/launch-edith.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/randywilson/Desktop/edith-v3</string>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/randywilson/Desktop/edith-v3/logs/edith.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/randywilson/Desktop/edith-v3/logs/edith.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/Users/randywilson/.bun/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/randywilson</string>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
