# packages/agent/com.edith.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.edith.agent</string>

    <key>Comment</key>
    <string>Edith AI Assistant — auto-starts on login</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/randywilson/.bun/bin/bun</string>
      <string>/Users/randywilson/Desktop/edith-v3/packages/agent/edith.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/randywilson/Desktop/edith-v3/packages/agent</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/randywilson/.edith/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/randywilson/.edith/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/randywilson</string>
      <key>PATH</key>
      <string>/Users/randywilson/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
